# netbirdui.service
[Unit]
Description=Netbird UI

[Service]
Type=simple
ExecStart=/usr/bin/netbird-ui

[Install]
WantedBy=graphical-session.target
